# com.paulbuckley.email-checker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.paulbuckley.email-checker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/paulbuckley/projects/work-tools/email-checker/venv/bin/python3</string>
        <string>/Users/paulbuckley/projects/work-tools/email-checker/app.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/paulbuckley/projects/work-tools/email-checker/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/paulbuckley/projects/work-tools/email-checker/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
